# example/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare message

## Get the 'message' key from the user config options.
message=$(bashio::config 'message')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${message:="Hello World..."}"

## Run your program
exec TF_LABELS=con,sin TF_LABEL_INVOKE=con MODEL_FILE=/addon_configs/model.tflite WATCHDIR=/config/camterraza WEBHOOKURL=http://192.168.0.58:8123/api/webhook/-zOxVfLI6XOYf6sAL2HZ0pd9Z python3 /usr/bin/my_program